# ansible/roles/wireguard/templates/wireguard.service.j2
[Unit]
Description=wireguard vpn
After=network.target

[Service]
Type=forking
ExecStart=/etc/wireguard/serviceup.sh

[Install]
WantedBy=multi-user.target
