# roles/snapraid/templates/snapraid-runner.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=SnapRAID runner

[Service]
Type=oneshot
User=root
Group=root
ExecStart={{ snapraid_runner_dir }}/run.sh
SyslogIdentifier=snapraid-runner
